# .config/systemd/user/package-export.timer
[Unit]
Description=Export a list of installed packages weekly

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
